# supervisor-sdns.ini
[program:dis-sdns]
command=/root/dis-sdns/output/control.sh start
directory=/root/dis-sdns/output
user=root
autostart=false
startsecs=3
startretries=3
autorestart=true
stdout_logfile=/var/log/dis-sdns.out
stderr_logfile=/var/log/dis-sdns.err